# systemd/neurocards-worker@.service
[Unit]
Description=Neurocards Video Generation Worker (Instance %i)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/neurocards/neurocards-bot
Environment="PYTHONUNBUFFERED=1"
Environment="WORKER_ID=%i"
EnvironmentFile=/var/neurocards/neurocards-bot/.env
ExecStart=/usr/bin/python3 -m worker.worker
Restart=always
RestartSec=10
StandardOutput=append:/var/neurocards/neurocards-bot/worker-%i.log
StandardError=append:/var/neurocards/neurocards-bot/worker-%i.log

[Install]
WantedBy=multi-user.target
